# supervisord.conf
[supervisord]
nodaemon=true

[program:ircd]
command=/opt/charybdis/bin/ircd -foreground
user=ircd
